# puppet/modules/cronjobs/templates/etc/cron.d/journeymonitor-infra.erb
# managed by puppet
MAILTO=""

<% if @env == "prod" %>
* * * * * root /opt/simplecd/simplecd.sh tag travisci-build-master-* https://github.com/journeymonitor/infra.git https://github.com/journeymonitor/infra/commit/
<% end %>
