# puppet/modules/cronjobs/templates/etc/cron.d/journeymonitor-analyze.erb
# managed by puppet
MAILTO=""

<% if @env == "prod" %>
* * * * * root /opt/simplecd/simplecd.sh tag travisci-build-master-* https://github.com/journeymonitor/analyze.git https://github.com/journeymonitor/analyze/commit/ 2>&1 | while IFS= read -r line;do echo "$(date) $line";done >> /var/tmp/simplecd/journeymonitor-analyze.log
<% end %>

*/5 * * * * journeymonitor /usr/bin/curl <%= scope.function_hiera(['app-analyze::endpoint_monitor_api']) %> -o <%= scope.function_hiera(['app-analyze::testresults_cache_path']) %> 2>&1 | while IFS= read -r line;do echo "$(date) $line";done >> /var/tmp/journeymonitor-analyze-import-and-spark.log && /bin/bash -c '. /etc/journeymonitor/app-analyze-env.sh && /usr/bin/java -jar /opt/journeymonitor/analyze/journeymonitor-analyze-importer-assembly.jar <%= scope.function_hiera(['app-analyze::testresults_cache_path']) %> 2>&1 | while IFS= read -r line;do echo "$(date) $line";done >> /var/tmp/journeymonitor-analyze-import-and-spark.log' && /bin/bash -c '. /etc/journeymonitor/app-analyze-env.sh && /opt/spark-1.5.1-bin-hadoop-2.6_scala-2.11/bin/spark-submit --deploy-mode client --executor-memory 256m --master spark://<%= scope.function_hiera(['app-analyze::spark_master_address']) %>:7077 /opt/journeymonitor/analyze/journeymonitor-analyze-spark-assembly.jar 2>&1 | while IFS= read -r line;do echo "$(date) $line";done >> /var/tmp/journeymonitor-analyze-import-and-spark.log'
