# puppet/modules/cronjobs/templates/etc/cron.d/journeymonitor-infra.erb
# managed by puppet
MAILTO=""

<% if @disable_infra_ci == true %>
# INFRA Continuous Delivery disabled
<% end %>

<% if @env == "prod" && @disable_infra_ci == false %>
* * * * * root /opt/simplecd/simplecd.sh tag travisci-build-master-* https://github.com/journeymonitor/infra.git https://github.com/journeymonitor/infra/commit/ 2>&1 | while IFS= read -r line;do echo "$(date) $line";done >> /var/tmp/simplecd/journeymonitor-infra.log
<% end %>
