# puppet/modules/cronjobs/templates/etc/cron.d/journeymonitor-other-spark-slave-cleanup.erb
# managed by puppet
MAILTO=""

# Remove Spark work data older than 1 hour
4 * * * * root /usr/bin/find /opt/spark-2.0.2-bin-hadoop2.7/work/* -maxdepth 0 -type d -mmin +59 | /usr/bin/xargs /bin/rm -rf
