# puppet/modules/cronjobs/templates/etc/cron.d/journeymonitor-control.erb
# managed by puppet
MAILTO=""

<% if @env == "prod" %>
* * * * * root /opt/simplecd/simplecd.sh tag travisci-build-master-* https://github.com/journeymonitor/control.git https://github.com/journeymonitor/control/commit/ 2>&1 | while IFS= read -r line;do echo "$(date) $line";done >> /var/tmp/simplecd/journeymonitor-control.log
<% end %>

* * * * * root sudo -u www-data /usr/bin/php /opt/journeymonitor/control/php/bin/console journeymonitor:control:import:testresults <%= scope.function_hiera(['cronjobs::endpoint_monitor_api']) %> --env=<%= @env %> 2>&1 | while IFS= read -r line;do echo "$(date) $line";done >> /var/tmp/journeymonitor-control-import-testresults.log

* * * * * root sudo -u www-data /bin/bash /opt/journeymonitor/control/statisticsimporter/start.sh 2>&1 | while IFS= read -r line;do echo "$(date) $line";done >> /var/tmp/journeymonitor-control-statisticsimporter.log
