# puppet/modules/cronjobs/templates/etc/cron.d/selenior-infra.erb
# managed by puppet
MAILTO=""

<% if @env == "prod" %>
* * * * * root /opt/simplecd/simplecd.sh branch master https://JourneyMonitorDeployer:Uhnd736gvfs5198G@bitbucket.org/selenior/infra.git
<% end %>
